# blacklist-hunter.cron
*/2 * * * * sh /use/local/bin/blacklist-hunter.sh
* * * * * root echo hello > /proc/1/fd/1 2>/proc/1/fd/2
